# templates/radarr.service.epp
## FILE  MANGED BY PUPPET
# DO NOT EDIT

[Unit]
Description=Radarr Daemon
After=syslog.target network.target

[Service]
User=radarr
Group=radarr

Type=simple

ExecStart=/usr/bin/mono --debug <%= $path %>Radarr.exe -nobrowser -data=<<%= $path %>data>
TimeoutStopSec=20
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target